# debian/dnscrypt-wrapper.init
#!/bin/sh
# Start/stop the DNSCrypt wrapper.
#
### BEGIN INIT INFO
# Provides