# debian/dnscrypt-wrapper.init
#!/bin/sh
# Start/stop the DNSCrypt wrapper.
#
### BEGIN INIT INFO
# Provides:          dnscrypt-wrapper
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $